# roles/xui/templates/x-ui.service.j2
# Managed by Ansible — do not edit manually.
[Unit]
Description=x-ui Service
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory={{ xui_install_dir }}/
ExecStart={{ xui_binary }}
Restart=on-failure
RestartSec=5s
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
